# deploy/fapai-crawler.timer
[Unit]
Description=法拍者联盟 Daily Crawler Timer (2 AM)

[Timer]
OnCalendar=*-*-* 02:00:00
Persistent=true
RandomizedDelaySec=600

[Install]
WantedBy=timers.target
