# deploy/fapai-crawler.service.d/timeout.conf
[Service]
# 爬虫正常几小时跑完;超过4小时视为卡死(如Playwright管道崩溃),自动终止避免空耗。
# 2026-06-11 修复:6/10曾因EPIPE卡死20小时。
TimeoutStartSec=14400
